# mapeditor.service
[Unit]
Description=OpenMower Map Editor Service
After=network.target

[Service]
# Radera ev. befintlig container utan att felet bubbla upp
ExecStartPre=-/usr/bin/podman rm -f openmower-mapeditor

# Kör containern i förgrunden (host network + volymer)
ExecStart=/usr/bin/podman run \
  --name openmower-mapeditor \
  --rm \
  --network host \
  -v /root/ros_home/.ros/map.bag:/data/map.bag:Z \
  -v /opt/openmower-mapeditor/data:/data:Z \
  openmower-mapeditor:latest

# Efter start, försök extrahera bag → geojson (ignorera fel)
ExecStartPost=-/usr/bin/podman exec openmower-mapeditor \
  rosbag2geojson /data/map.bag /data/map.geojson

# Se till att containern städas upp på stopp
ExecStopPost=-/usr/bin/podman rm -f openmower-mapeditor

Restart=on-failure
RestartSec=5s
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
